# real-money-test/systemd/freqtrade-macd-aggressive-dryrun.service.example
[Unit]
Description=Freqtrade MACD Aggressive Dry Run
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/path/to/quant/test3
ExecStart=/path/to/quant/test3/real-money-test/start_dry_run.sh
ExecStop=/bin/kill -s SIGINT $MAINPID
Restart=on-failure
RestartSec=10
TimeoutStopSec=120
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
